# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:apache]
command=apache2-foreground
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:scheduler]
# Laravel's schedule:run is meant to be invoked every minute by cron. Looping with
# sleep 60 gives the same effect inside a single container. This is where booking
# reminders, slot-release jobs, and tournament progression will eventually hang off.
command=/bin/bash -c "while true; do php artisan schedule:run --no-interaction; sleep 60; done"
directory=/var/www/html
# Run as www-data (same user as Apache) so any cache/session files written by
# scheduled commands stay writable by the web worker.
user=www-data
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
